# config/deploy/templates/shared/foreman/process.service.erb
[Unit]
PartOf=<%= app %>-<%= name %>.target

[Service]
# Foreground process (do not use --daemon in ExecStart or config.rb)
Type=simple
WorkingDirectory=<%= engine.root %>
ExecStart=/bin/bash -l -c '<%= process.command %>'
Restart=no
TimeoutStopSec=<%= engine.options[:timeout] %>
SyslogIdentifier=deploy-it.<%= name %>
